# root/etc/s6-overlay/s6-rc.d/init-zt/run
#!/command/with-contenv bash

if [ ! -f "/var/lib/zerotier-one/zerotier-cli" ]; then
  cp -r /var/lib/zerotier-one.bak/ /var/lib/zerotier-one/
fi