# root/etc/s6-overlay/s6-rc.d/init-warp-ip-selection/run
#!/command/with-contenv bash
set -uo pipefail

WARP_IP_SELECTION_ENABLED="${WARP_IP_SELECTION_ENABLED:-false}"
WARP_API_SELECTION_ENABLED="${WARP_API_SELECTION_ENABLED:-false}"
WARP_OVERRIDE_WARP_ENDPOINT="${WARP_OVERRIDE_WARP_ENDPOINT:-}"
WARP_OVERRIDE_API_ENDPOINT="${WARP_OVERRIDE_API_ENDPOINT:-}"

if [ "$WARP_IP_SELECTION_ENABLED" != "true" ] && [ "$WARP_API_SELECTION_ENABLED" != "true" ]; then
  echo "[ip-selection] Endpoint selection disabled"
  exit 0
fi

# ── API 优选 ──
# API 优选结果仅写入 WARP_OVERRIDE_API_ENDPOINT，不再覆盖 Tunnel endpoint
if [ "$WARP_API_SELECTION_ENABLED" = "true" ]; then
  if [ -n "$WARP_OVERRIDE_API_ENDPOINT" ]; then
    echo "[ip-selection] API endpoint already provided: ${WARP_OVERRIDE_API_ENDPOINT}"
  else
    BEST_API_ENDPOINT="$(warp-speed-test.sh --api 2>/dev/null || true)"
    BEST_API_ENDPOINT="$(echo "${BEST_API_ENDPOINT}" | tr -d '[:space:]')"
    if [ -n "$BEST_API_ENDPOINT" ]; then
      echo "$BEST_API_ENDPOINT" > /var/run/s6/container_environment/WARP_OVERRIDE_API_ENDPOINT
      echo "[ip-selection] API endpoint selected: ${BEST_API_ENDPOINT}"
    else
      echo "[ip-selection] API selection soft-failed"
    fi
  fi
fi

# ── Tunnel 优选 (独立运行，不受 API 结果影响) ──
if [ "$WARP_IP_SELECTION_ENABLED" = "true" ]; then
  if [ -n "$WARP_OVERRIDE_WARP_ENDPOINT" ]; then
    echo "[ip-selection] Tunnel endpoint already provided: ${WARP_OVERRIDE_WARP_ENDPOINT}"
  else
    BEST_TUNNEL_ENDPOINT="$(warp-speed-test.sh --tunnel 2>/dev/null || true)"
    BEST_TUNNEL_ENDPOINT="$(echo "${BEST_TUNNEL_ENDPOINT}" | tr -d '[:space:]')"
    if [ -n "$BEST_TUNNEL_ENDPOINT" ]; then
      echo "$BEST_TUNNEL_ENDPOINT" > /var/run/s6/container_environment/WARP_OVERRIDE_WARP_ENDPOINT
      echo "[ip-selection] Tunnel endpoint selected: ${BEST_TUNNEL_ENDPOINT}"
    else
      echo "[ip-selection] Tunnel selection soft-failed, keep default endpoint"
    fi
  fi
fi

exit 0
