# root/etc/s6-overlay/s6-rc.d/svc-warp/run
#!/command/with-contenv bash
exec /usr/bin/warp-svc
