# root/etc/s6-overlay/s6-rc.d/svc-dbus/run
#!/command/with-contenv bash
exec /usr/bin/dbus-daemon --system --nofork
